# manjaro/xfce/desktop-overlay/etc/skel/.config/autostart/custom-dock.desktop
[Desktop Entry]
Name=My Script
Exec=/usr/share/custom-dock/center-dock.sh
Type=Application
