# manjaro/xfce/desktop-overlay/etc/skel/.config/autostart/custom-resolution-watch.desktop
[Desktop Entry]
Name=Resolution Watch
Exec=sh ~/.config/custom-dock/resolution-watch.sh
Type=Application
