# manjaro/xfce/desktop-overlay/etc/skel/.config/autostart/custom-dock.desktop
[Desktop Entry]
Name=Center-dock
Exec=sh -c '~/.config/custom-dock/dock-position.sh'
Type=Application
